# dot_config/systemd/user/bt-a2dp-monitor.service
[Unit]
Description=Monitor Bluetooth connection events for A2DP recovery
After=pipewire.service wireplumber.service
ConditionPathExists=%h/.config/bt-a2dp-fix/config
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
ExecStart=%h/.local/bin/bt-a2dp-monitor.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
